# conf/etc/periodic/15min/01-drupal
#!/bin/sh
cd /var/www/drupal
vendor/bin/drush cron
